# openclaw-manager.service
[Unit]
Description=OpenClaw Manager
After=network.target

[Service]
Type=simple
User=chenyi
WorkingDirectory=/home/chenyi/.openclaw/workspace/projects/openclaw-manager
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=3
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
